# docker/alpine/etc/wingmate/crontab
*/5 * * * * /etc/wingmate/crontab.d/cron1.sh